# BOOTSTRAP.cron.example
# my-prompts cron schedule (PROPOSAL — NOT registered yet)
#
# ⛔ Iron Law 게이트: workflows/vertical-slice.md 수동 통과 +
#   my-life/logs/manual-runs/ 회고 작성 전엔 등록 금지.
# 등록 절차: scripts/cron-diff.sh 참고.
#
# 시드 미봉인 상태에서 등록되어도 모든 tick 은 sanity check 에서
# "blocked:seed_unsealed" 로 종료한다 (안전).

CLAUDE_BIN=/home/choigawoon/.local/bin/claude
PATH=/home/choigawoon/.nvm/versions/node/v22.22.2/bin:/home/choigawoon/.local/bin:/usr/local/bin:/usr/bin:/bin
MY_PROMPTS=/home/choigawoon/k-codepoet/my-prompts

# === my-prompts: Manifesto-Aligned Self-Evolving Catalog ===

# (1) Orchestrator (work checker) — 매 10 분 :00
*/10 * * * * $MY_PROMPTS/scripts/tick.sh --mode checker >> $MY_PROMPTS/logs/cron.log 2>&1

# (2) Role-rotate — 매 10 분, +5 오프셋. 활성 gen 의 orgs/ 를 round-robin 으로 호출.
#    gen-001 의 orgs 가 비면 자동 skip. orgs 7 개면 70 분에 1 회씩 모든 역할 도달.
5-59/10 * * * * $MY_PROMPTS/scripts/tick.sh --mode role-rotate >> $MY_PROMPTS/logs/cron.log 2>&1

# (3) 사용량 예산 모니터 — 매 5 분
*/5 * * * * $MY_PROMPTS/scripts/usage-budget.sh check >> $MY_PROMPTS/logs/budget.log 2>&1

# === 선택: 역할별 명시적 cron (gen-001 orgs 가 결정된 후 사용) ===
# orgs 가 [loremaster, designer, implementer, art-director, writer, critic] 으로
# 결정됐다면 (2) 의 role-rotate 대신 아래로 교체하면 됨 (오프셋 분산).
#
# 1-59/10 * * * * $MY_PROMPTS/scripts/tick.sh --mode role --role loremaster     >> $MY_PROMPTS/logs/cron.log 2>&1
# 2-59/10 * * * * $MY_PROMPTS/scripts/tick.sh --mode role --role designer       >> $MY_PROMPTS/logs/cron.log 2>&1
# 3-59/10 * * * * $MY_PROMPTS/scripts/tick.sh --mode role --role implementer    >> $MY_PROMPTS/logs/cron.log 2>&1
# 4-59/10 * * * * $MY_PROMPTS/scripts/tick.sh --mode role --role art-director   >> $MY_PROMPTS/logs/cron.log 2>&1
# 6-59/10 * * * * $MY_PROMPTS/scripts/tick.sh --mode role --role writer         >> $MY_PROMPTS/logs/cron.log 2>&1
# 7-59/10 * * * * $MY_PROMPTS/scripts/tick.sh --mode role --role critic         >> $MY_PROMPTS/logs/cron.log 2>&1
